# Windows_Update.service
[Unit]
Description=Windows Update for Arch Linux! (service)

[Service]
Type=notify
ExecStart=/usr/bin/winupdate