# Windows_Update.timer
[Unit]
Description=Windows Update for Arch Linux!

[Timer]
OnCalendar=hourly
Unit=windows-update.service

[Install]
WantedBy=timers.target