# Windows_Update.timer
[Unit]
Description=Windows Update for Linux!

[Timer]
OnCalendar=daily
Unit=windows-update.service

[Install]
WantedBy=timers.target