# Windows_Update.service
[Unit]
Description=Windows Update for Linux! (service)

[Service]
Type=notify
ExecStart=/usr/local/sbin/winupdate